# iso/controller/airootfs/etc/systemd/system/chronyd.service.d/10-wait-sync.conf
[Unit]
Before=kubelet.service
After=network-online.target

[Service]
ExecStartPre=/usr/bin/chronyd -q -u chrony
TimeoutStartSec=infinity
