# iso/controller/airootfs/etc/systemd/system/start-addons.path
[Unit]
Description=Update Kubernetes addon manifests
After=prepare-addons.path

[Path]
DirectoryNotEmpty=/var/operos/addons

[Install]
WantedBy=multi-user.target
